# scripts/com.viren.nas-auto-cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.viren.nas-auto-cleanup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/viren/scripts/nas-auto-cleanup.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>3600</integer> <!-- run every hour -->
  
  <key>StandardOutPath</key>
  <string>/tmp/nas-auto-cleanup.out</string>
  
  <key>StandardErrorPath</key>
  <string>/tmp/nas-auto-cleanup.err</string>
</dict>
</plist>
